# etc/supervisor/conf.d/redis-server.conf
[program:redis-server]
command=/usr/bin/redis-server
stdout_logfile=%(ENV_VIRTUAL_ENV)s/var/log/redis-server/out.log 
stderr_logfile=%(ENV_VIRTUAL_ENV)s/var/log/redis-server/err.log
stopasgroup=true
autostart=true
autorestart=false
