# etc/supervisor/conf.d/tlsensing-server.conf
[program:tlsensing-server]
command=%(ENV_VIRTUAL_ENV)s/TLSensing/inits/init_tls_server.sh
stdout_logfile=%(ENV_VIRTUAL_ENV)s/var/log/tlsensing-server/out.log 
stderr_logfile=%(ENV_VIRTUAL_ENV)s/var/log/tlsensing-server/err.log
stopasgroup=true
autostart=true
autorestart=false
