# TLSensing/bootstrap/data/etc/supervisor/conf.d/acquire-values-oneshot.conf
[program:acquire-values-oneshot]
command=%(ENV_VIRTUAL_ENV)s/TLSensing/inits/acquire-values-oneshot.sh
stdout_logfile=%(ENV_VIRTUAL_ENV)s/var/log/acquire-values/out.log
stderr_logfile=%(ENV_VIRTUAL_ENV)s/var/log/acquire-values/err.log
stopasgroup=true
autostart=false
autorestart=false
